# Pi/zram/zram.service
[Unit]
Description=zram
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/zram.sh

[Install]
WantedBy=multi-user.target
